# templates/pods/caddy/caddy.container.j2
[Unit]
Description=Caddy Quadlet

[Container]
Network=services.network
ContainerName=caddy-reverse-proxy
Label=app=caddy
AutoUpdate=registry
Image=docker.io/library/caddy:2.10
Volume=caddy-data.volume:/data
Volume=caddy-config.volume:/config
Volume=%h/Caddyfile:/etc/caddy/Caddyfile:Z
# Volume=website.volume:/var/www/html
# static IP for nextcloud to trust
IP=10.89.2.58
PublishPort={{ caddy_http_port }}:80
PublishPort={{ caddy_https_port }}:443

[Service]
Restart=on-failure
TimeoutStartSec=900

[Install]
WantedBy=default.target
