# templates/pods/paperless/paperless-webserver.container
[Container]
Pod=paperless.pod
ContainerName=paperless_ngx
Label=app=paperless
AutoUpdate=registry
Image=ghcr.io/paperless-ngx/paperless-ngx:2.18
Environment=PAPERLESS_REDIS=redis://localhost:6379 PAPERLESS_DBHOST=localhost
Environment=USERMAP_UID=0 USERMAP_GID=0
# search index, SQLite database, classification model, etc
Environment=PAPERLESS_DATA_DIR=/opt/paperless/data
# user-uploaded documents and thumbnails
Environment=PAPERLESS_MEDIA_ROOT=/opt/paperless/media
# alternative location to upload documents to paperless
Environment=PAPERLESS_CONSUMER_DISABLE=true
{% if podman_services.paperless.enable_tika is sameas true %}
Environment=PAPERLESS_TIKA_ENABLED=1 PAPERLESS_TIKA_GOTENBERG_ENDPOINT=http://localhost:3000 PAPERLESS_TIKA_ENDPOINT=http://localhost:9998
{% endif %}
# Environment=PAPERLESS_CONSUMPTION_DIR=/opt/paperless/consume
Secret=paperless-db-password,type=env,target=PAPERLESS_DBPASS
# the data directories specified above are auto-created by paperless container,
# so we just need to create /data/paperless on host
Volume=/data/paperless:/opt/paperless:Z

[Unit]
Requires=paperless-db.service paperless-broker.service
After=paperless-db.service paperless-broker.service
