# templates/pods/paperless/paperless-tika.container
[Container]
Pod=paperless.pod
ContainerName=paperless_tika
Label=app=paperless
Image=docker.io/apache/tika:3.2.3.0

