# templates/restic-backup.service.j2
[Unit]
Description=Restic Backup Service

[Service]
Type=oneshot
User={{ ansible_user }}
Environment="RESTIC_PASSWORD={{ restic_backup_password }}"
ExecStart=/usr/local/bin/backup.sh
