# templates/pods/immich/immich-server.container
[Container]
Pod=immich.pod
ContainerName=immich_server
Label=app=immich
AutoUpdate=registry
Image=ghcr.io/immich-app/immich-server:release
Environment=DB_HOSTNAME=localhost REDIS_HOSTNAME=localhost
Secret=immich-db-password,type=env,target=DB_PASSWORD
# thumbnails and database dumps are also placed here
# https://immich.app/docs/administration/backup-and-restore
Volume=/data/photos:/data
Volume=/etc/localtime:/etc/localtime:ro
# for quicksync
AddDevice=/dev/dri:/dev/dri

[Unit]
After=immich-database.service immich-redis.service
