# templates/pods/immich/immich-database.container
[Container]
Pod=immich.pod
ContainerName=immich_postgres
Label=app=immich
Image=ghcr.io/immich-app/postgres:14-vectorchord0.4.3
Notify=true
Environment=POSTGRES_DB=immich POSTGRES_INITDB_ARGS=--data-checksums
Secret=immich-db-password,type=env,target=POSTGRES_PASSWORD
Volume=immich-db.volume:/var/lib/postgresql/data
ShmSize=128mb

[Service]
Restart=always
