# templates/pods/miniflux/miniflux.container
[Container]
Pod=miniflux.pod
ContainerName=miniflux
Label=app=miniflux
# in general, auto-updates should be fine with RUN_MIGRATIONS: https://miniflux.app/docs/upgrade.html
# but there is only a `latest` tag, which I would rather not blindly follow
# AutoUpdate=registry
Image=ghcr.io/miniflux/miniflux:2.2.12
Environment=BASE_URL=https://feeds.genevievemendoza.com DATABASE_URL=postgres://miniflux:secret@localhost/miniflux?sslmode=disable RUN_MIGRATIONS=1 POLLING_FREQUENCY=300
# on the first run you must also specify env vars:
# CREATE_ADMIN=1 ADMIN_USERNAME=admin ADMIN_PASSWORD=test123
# then log in and change the credentials yourself

[Unit]
After=miniflux-pg.service
